# luci-app-tlsfragment/root/etc/uci-defaults/40_luci-tlsfragment
#!/bin/sh

uci -q batch <<-EOF >/dev/null
    delete ucitrack.@tlsfragment[-1]
    add ucitrack tlsfragment
    set ucitrack.@tlsfragment[-1].init=tlsfragment
    commit ucitrack
